# files/etc/service/seahub/finish
#!/bin/sh

/opt/seafile/seafile-server-latest/seahub.sh stop >>/var/log/seafile/seahub.log 2>&1
